# infra/coreos/cockpit.container
[Install]
WantedBy=default.target

[Container]
Image=ghcr.io/arlyon/cockpit:latest
Network=network.network
